# public/generated/posts/blog/85.json
{
  "day": 85,
  "title": "Understanding Cryptomining Malware",
  "date": "2026-04-26",
  "url": "https://juribuora.github.io/blog/2026/04/26/day-85.html",
  "category": "blog",
  "tags": [
    "cybersecurity",
    "processes",
    "learningprocess"
  ],
  "slug": "day-85",
  "sourcePath": "Blog/_posts/2026-04-26-day-85.md",
  "content": "# 🔄 Topic\n\nUnderstanding cryptomining malware and how attackers abuse compromised systems for profit.\n\n---\n\n# 🎯 Goal\n\nLearn how cryptomining malware works, what indicators it creates, and why it can remain hidden longer than more destructive malware.\n\n---\n\n# 🛠 What I Did\n\nToday I studied cryptomining malware.\n\nCryptomining malware uses compromised systems to mine cryptocurrency for the attacker.\n\nUnlike ransomware, it does not usually encrypt files or immediately stop business operations.\n\nInstead, it quietly abuses system resources.\n\nCommon indicators include:\n\n- high CPU usage\n- high GPU usage\n- unusual PowerShell activity\n- unknown processes consuming resources\n- connections to mining pools\n- systems becoming slow\n- unexpected scheduled tasks\n- persistence mechanisms\n\nOne example I reviewed was WannaMine, a cryptomining malware family that abused Windows systems to mine cryptocurrency.\n\nThe important distinction:\n\nRansomware wants fast impact.\n\nCryptomining often wants long-term resource theft.\n\n---\n\n# 🔗 Key Cybersecurity Connections\n\nCryptomining matters because it may be overlooked.\n\nIf a system is slow, users may assume:\n\n- old hardware\n- too many applications open\n- normal workload\n- bad performance day\n\nBut from a security perspective, persistent high resource usage can be evidence.\n\nCryptomining malware may enter through:\n\n- vulnerable services\n- stolen credentials\n- malicious scripts\n- exposed cloud workloads\n- weak container security\n- unpatched systems\n\nIn cloud environments, cryptomining can also create major financial cost because attackers consume compute resources.\n\nSo the impact is not only performance.\n\nIt can also be money.\n\n---\n\n## 🔍 Investigation Questions\n\n- Which process is consuming CPU or GPU?\n- Is the process expected?\n- Who launched it?\n- Where is the executable located?\n- Is PowerShell involved?\n- Is the host connecting to mining pools?\n- Are there suspicious scheduled tasks?\n- Is the activity persistent after reboot?\n- Did the host recently connect to suspicious domains?\n- Are multiple hosts showing the same behavior?\n- Is this happening on servers, endpoints, or cloud workloads?\n\n---\n\n## 🚨 Detection Opportunities\n\nPossible detection ideas:\n\n- sustained high CPU usage by unknown process\n- endpoint connecting to known mining pool domains\n- unusual process names imitating system processes\n- PowerShell downloading mining payloads\n- scheduled task executing suspicious binary\n- container using abnormal CPU\n- cloud instance with unexpected compute spike\n- outbound traffic to mining-related infrastructure\n- many endpoints connecting to the same mining pool\n\nExample detection pattern:\n\n    process=unknown.exe\n    cpu_usage=95%\n    duration=45 minutes\n    destination_domain=mining-pool.example\n    suspicion=cryptomining_malware\n\nThis combines endpoint and network evidence.\n\nThat is stronger than CPU usage alone.\n\n---\n\n## 🧭 MITRE ATT&CK Techniques\n\n- T1496 — Resource Hijacking\n- T1059 — Command and Scripting Interpreter\n- T1105 — Ingress Tool Transfer\n- T1053 — Scheduled Task/Job\n- T1071 — Application Layer Protocol\n\n---\n\n# 🗺 Visual Investigation Diagram\n\n    Initial compromise\n        ↓\n    Miner payload delivered\n        ↓\n    Persistence created\n        ↓\n    CPU/GPU resources abused\n        ↓\n    Mining pool contacted\n        ↓\n    Attacker profits quietly\n\n---\n\n# ⚠ Challenges\n\nThe main challenge is that cryptomining may not create obvious damage.\n\nThe system still runs.\n\nThe user may only notice slowness.\n\nThat makes the infection easier to ignore.\n\nAnother challenge is that performance problems are not automatically security problems.\n\nHigh CPU could be legitimate.\n\nThe defender needs to combine:\n\n- process evidence\n- user context\n- network destinations\n- persistence artifacts\n- historical baseline\n\n---\n\n# 📚 What I Learned\n\nI learned that not all malware wants to destroy or encrypt data.\n\nSome malware wants to stay quiet and profit over time.\n\nThis changes the detection mindset.\n\nInstead of only looking for dramatic impact, I need to look for resource abuse, persistence, and suspicious outbound communication.\n\n---\n\n# ➡ Next Steps\n\n- Study mining pool traffic patterns\n- Review endpoint process monitoring\n- Practice investigating high CPU alerts\n- Learn cloud cryptomining detection patterns\n- Compare cryptomining malware with botnet behavior\n- Add resource hijacking to my threat notes\n\n---\n\n# 🧠 Reflection\n\nCryptomining malware is interesting because it is quieter than ransomware but still harmful.\n\nIt shows that attacker motivation changes attacker behavior.\n\nSome attackers want disruption.\n\nSome want access.\n\nSome want compute power.\n\nThe SOC analyst needs to recognize all of them.\n\n---\n\n# 🧩 Lessons Learned\n\n**What worked**\n\nConnecting performance indicators with security investigation.\n\n**What broke**\n\nThinking malware always causes immediate visible damage.\n\n**Why it broke**\n\nCryptomining can profit from slow, quiet resource abuse.\n\n**Fix / takeaway**\n\nSustained abnormal resource usage deserves investigation, especially when combined with suspicious network traffic.\n\n---\n\n# 📈 Skill Progression Context\n\nThis supports SOC readiness because cryptomining alerts are common in endpoint, cloud, and container environments.\n\nInvestigating them requires correlating process telemetry, network activity, and resource usage.\n\n---\n\n# 😄 TL;DR\n\nRansomware kicks the door down.\n\nCryptomining sneaks in and uses your electricity bill.\n"
}
